# devops/roles/icos.docker/files/docker-periodic-cleanup.timer
[Unit]
Description=Cleanup dangling images once a week

[Timer]
# Run once a week, at midnight
OnCalendar=weekly

[Install]
WantedBy=timers.target
