# devops/roles/icos.doi/templates/doi.service.j2
# /etc/systemd/system/doi.service
[Unit]
Description=DOI minting service
After=network.target

[Service]
Type=simple
User={{ username }}
ExecStart=/usr/bin/java -jar "{{ jarservice_jar }}"
Restart=always
RestartSec=30
WorkingDirectory=/home/{{ username }}

[Install]
WantedBy=multi-user.target
